# crontab
# crontab
0 * * * * cd /app && python -c "from app.services.collector import main; main()" >> /var/log/collector.log 2>&1

