# crontab
# crontab
0 * * * * cd /app && python /app/run_collector.py >> /var/log/collector.log 2>&1

