# crontab
# crontab
0 * * * * cd /app && python -m app.services.collector >> /var/log/collector.log 2>&1

